# .docker/crontab
* * * * * root /usr/local/bin/php /var/www/artisan schedule:run >> /var/www/scheduler/cronjob.log 2>&1
